# modules/sinatraservice/files/sinatra.service
[Unit]
Description = Sinatra
After = network.target
[Service]
ExecStart = /apps/simple-sinatra-app/sinatra.sh
[Install]
WantedBy = multi-user.target